# services/nomad/root/etc/services.d/nomad/run
#!/usr/bin/with-contenv sh

# exec /usr/local/bin/nomad agent -dev -bind="0.0.0.0" -log-level=INFO -config=/etc/nomad/nomad.json;
